# scripts/nidas_eol_relay/eol_relay.crontab
*/10 * * * *    $HOME/nidas_eol_relay/check_nidas_udp_relay.sh start
