# scripts/nidas_eol_relay/eol_relay.crontab
PATH=/usr/bin:/opt/nidas/bin
*/10 * * * *    $HOME/check_nidas_udp_relay.sh
